# templates/dewy.service.erb
[Unit]
Description=The application server for deployment with distributed polling.
Documentation=https://github.com/linyows/dewy
Requires=network-online.target
After=network-online.target

[Service]
User=<%= $user %>
Group=<%= $group %>
Environment=GOMAXPROCS=2
Restart=on-failure
WorkingDirectory=<%= $dir %>
EnvironmentFile=<%= $sysconfig_dir %>/dewy
ExecStart=<%= $bin %>/dewy <%= $command %>
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target